# deploy/mcp-server/at-directory-mcp.service
[Unit]
Description=AT Directory MCP server (hosted HTTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/at-directory-mcp
Environment=PORT=8099
# AT_VERIFIER_URL intentionally unset: reads are ungated, so the server runs
# the full anonymous discovery+transact path without it. Set it (and
# daemon-reload + restart) once the AT verifier route is live to enable
# credentialed write features.
# Node 20 is installed side-by-side at /opt/node20; system Node 18 (OP API
# dependency) is never touched.
ExecStart=/opt/node20/bin/node /opt/at-directory-mcp/http.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
